# resources/fty-sensor-gpio.service.in
[Unit]
Description=fty-sensor-gpio service
Requires=network.target malamute.service fty-asset.service fty-info.service
After=network.target malamute.service fty-asset.service fty-info.service
PartOf=bios.target
BindsTo=fty-asset.service

[Service]
Type=simple
User=bios
Group=gpio
EnvironmentFile=-@CMAKE_INSTALL_PREFIX@/share/bios/etc/default/bios
EnvironmentFile=-@CMAKE_INSTALL_PREFIX@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@CMAKE_INSTALL_PREFIX@/share/fty/etc/default/fty
EnvironmentFile=-@CMAKE_INSTALL_PREFIX@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@CMAKE_INSTALL_SYSCONFDIR@/default/bios
EnvironmentFile=-@CMAKE_INSTALL_SYSCONFDIR@/default/bios__%n.conf
EnvironmentFile=-@CMAKE_INSTALL_SYSCONFDIR@/default/fty
EnvironmentFile=-@CMAKE_INSTALL_SYSCONFDIR@/default/fty__%n.conf
Environment="prefix=@CMAKE_INSTALL_PREFIX@"
Environment='SYSTEMD_UNIT_FULLNAME=%n'
ExecStart=@prefix@/bin/fty-sensor-gpio -c @CMAKE_INSTALL_SYSCONFDIR@/@PROJECT_NAME@/@PROJECT_NAME@.cfg
# Try the usual default SIGTERM, but don't wait forever if service locked-up
TimeoutStopSec=5
SuccessExitStatus=SIGKILL SIGTERM 143
Restart=always

[Install]
WantedBy=bios.target
